# supervisor-gftp.conf
[program:gftp]
command=/usr/sbin/globus-gridftp-server -d ALL -L /var/log/ -no-detach -c /etc/gridftp.conf -banner "Welcome to gridftp server"
user=root
autorestart=true
